# roles/vdirsyncer/templates/vdirsyncer-clean-sync.timer.j2
[Unit]
Description=10-min vdirsyncer attendee-stripping sync

[Timer]
OnCalendar=*:0/10
Persistent=true
Unit=vdirsyncer-clean-sync.service

[Install]
WantedBy=timers.target
