# roles/vdirsyncer/templates/vdirsyncer-clean-sync.timer.j2
[Unit]
Description=Hourly vdirsyncer attendee-stripping sync

[Timer]
OnCalendar=hourly
Persistent=true
Unit=vdirsyncer-clean-sync.service

[Install]
WantedBy=timers.target
