# roles/vdirsyncer/templates/vdirsyncer-clean-sync.service.j2
[Unit]
Description=Sync work calendars (datashare + exchange + icloud_privat) into iCloud
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/vdirsyncer/env
User=vdirsyncer
ExecStartPre=/bin/sh -c 'test -d /var/lib/vdirsyncer/.vdirsyncer/status/datashare_to_local && test -d /var/lib/vdirsyncer/.vdirsyncer/status/exchange_to_local && test -d /var/lib/vdirsyncer/.vdirsyncer/status/icloud_privat_to_local && test -d /var/lib/vdirsyncer/.vdirsyncer/status/merged_to_icloud || yes | /usr/bin/vdirsyncer discover'
WorkingDirectory=/var/lib/vdirsyncer
CacheDirectory=vdirsyncer
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictNamespaces=true
RestrictSUIDSGID=true
RestrictRealtime=true
SystemCallArchitectures=native
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ReadOnlyPaths=/etc/vdirsyncer
ReadWritePaths=/var/lib/vdirsyncer
ExecStart=/usr/bin/vdirsyncer sync datashare_to_local
ExecStart=/usr/bin/vdirsyncer sync exchange_to_local
ExecStart=/usr/bin/vdirsyncer sync icloud_privat_to_local
ExecStartPost=/usr/local/bin/anonymize_and_merge.sh
ExecStartPost=-/usr/bin/vdirsyncer sync merged_to_icloud
